# docker/supervisor_redis.ini
[program:redis]
command=redis-server /etc/redis.conf --daemonize no
redirect_stderr=true